# services/sellmate-poller.service
[Unit]
Description=SellMate Pi Poller
After=network-online.target vend-api.service
Wants=network-online.target vend-api.service

[Service]
Type=simple
User=sambeals
WorkingDirectory=/home/sambeals/sellmate-pi/app
Environment=VEND_API_KEY=CHANGE_ME
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/sambeals/sellmate-pi/.venv/bin/python /home/sambeals/sellmate-pi/app/cloud_poller.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
